# debian/DEBIAN/prerm
#!/bin/bash
systemctl stop fimpui.service
systemctl disable fimpui.service
