# debian/DEBIAN/prerm
#!/bin/bash
systemctl stop thingsplex.service
systemctl disable thingsplex.service
